# local-test/test.service

[Unit]
Description=command

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/pwd

[Install]
WantedBy=local.target
